# s6/s6-rc.d/log-prep/up
#!/bin/execlineb -P

if { mkdir -p /var/log/s6 }
if { chown nobody:nogroup /var/log/s6 }
chmod 02755 /var/log/s6
