# s6/s6-rc.d/tailscaled-log/run
#!/bin/execlineb -P

exec logutil-service /var/log/s6/tailscale
